# etc/supervisord.d/mariadb.conf.ini
[program:mariadb]
command=/usr/local/bin/run_mariadb
