# etc/supervisord.d/mariadb.conf.ini
[program:mariadb]
command=/usr/bin/pidproxy /var/lib/mysql/mysqld.pid /usr/bin/mysqld_safe --transaction-isolation=%(ENV_TRANSACTION_ISOLATION_LEVEL)s
