# etc/supervisord.d/mariadb.conf.ini
[program:mariadb]
command=/usr/sbin/mysqld --transaction-isolation=%(ENV_TRANSACTION_ISOLATION_LEVEL)s %(ENV_MARIADB_OPTS)s --user=mysql
